# ai-stack.service
[Unit]
Description=AI Stack (Ollama, Whisper, Piper, Frigate, Postgres, LiteLLM, RouteLLM, nginx, Open WebUI, n8n, Monitoring)
Requires=docker.service
After=docker.service local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/sbin/nvidia-smi
ExecStart=/usr/bin/docker compose -f /ai/compose/ollama.yaml up -d
ExecStart=/usr/bin/docker compose -f /ai/compose/whisper.yaml up -d
ExecStart=/usr/bin/docker compose -f /ai/compose/piper.yaml up -d
ExecStart=/usr/bin/docker compose -f /ai/compose/frigate.yaml up -d
ExecStart=/usr/bin/docker compose -f /ai/compose/postgres.yaml up -d
ExecStart=/usr/bin/docker compose -f /ai/compose/litellm.yaml up -d
ExecStart=/usr/bin/docker compose -f /ai/compose/routellm.yaml up -d
ExecStart=/usr/bin/docker compose -f /ai/compose/openwebui.yaml up -d
ExecStart=/usr/bin/docker compose -f /ai/compose/n8n.yaml up -d
ExecStart=/usr/bin/docker compose -f /ai/compose/mcp.yaml up -d
ExecStart=/usr/bin/docker compose -f /ai/compose/nginx.yaml up -d
ExecStart=/usr/bin/docker compose -f /ai/compose/monitoring.yaml up -d
# Orchestrator is still in development — uncomment when ready:
#ExecStart=/usr/bin/docker compose -f /ai/compose/orchestrator.yaml up -d
ExecStop=/usr/bin/docker compose -f /ai/compose/monitoring.yaml down
ExecStop=/usr/bin/docker compose -f /ai/compose/orchestrator.yaml down
ExecStop=/usr/bin/docker compose -f /ai/compose/nginx.yaml down
ExecStop=/usr/bin/docker compose -f /ai/compose/mcp.yaml down
ExecStop=/usr/bin/docker compose -f /ai/compose/n8n.yaml down
ExecStop=/usr/bin/docker compose -f /ai/compose/openwebui.yaml down
ExecStop=/usr/bin/docker compose -f /ai/compose/routellm.yaml down
ExecStop=/usr/bin/docker compose -f /ai/compose/litellm.yaml down
ExecStop=/usr/bin/docker compose -f /ai/compose/postgres.yaml down
ExecStop=/usr/bin/docker compose -f /ai/compose/frigate.yaml down
ExecStop=/usr/bin/docker compose -f /ai/compose/piper.yaml down
ExecStop=/usr/bin/docker compose -f /ai/compose/whisper.yaml down
ExecStop=/usr/bin/docker compose -f /ai/compose/ollama.yaml down
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
